# piper-tts.service
[Unit]
Description=Nabaztag Piper TTS Proxy
Documentation=https://github.com/jsapede/nabaztag-piper
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=__INSTALL_DIR__
EnvironmentFile=__INSTALL_DIR__/.env
ExecStart=/usr/bin/python3 __INSTALL_DIR__/piper_tts_stream.py --phonemes
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nabaztag-piper

[Install]
WantedBy=multi-user.target
