# piper-tts.service
[Unit]
Description=Piper TTS Proxy for Nabaztag
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/New_Serverless_v3
EnvironmentFile=/opt/New_Serverless_v3/.env
ExecStart=/usr/bin/python3 /opt/New_Serverless_v3/piper_tts_stream.py --phonemes
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
